# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=python -u main.py

[program:healthcheck]
command=python -u healthcheck.py 